# tools/build/supremica_build.cron
#!/bin/bash
/usr/local/bin/supremica_automated_build >& /tmp/supremica_buildlog

exit 0


